# scanner.service
[Unit]
Description=Scanner Service

[Service]
ExecStart=/home/pi/scanner.py -d
Restart=always

[Install]
WantedBy=multi-user.target
